# rootfs/etc/s6-overlay/s6-rc.d/proxysql/run
#!/usr/bin/env bash

if [ -n "$CONFILEFILE" ]; then
    proxysql -f --idle-threads -D /var/lib/proxysql -c $CONFILEFILE
fi
proxysql -f --idle-threads -D /var/lib/proxysql
